# dotfiles/common/Scripts/systemd_scripts/chown-etc.service
[Unit]
Description=Chown etc Folder
After=multi-user.target

[Service]
ExecStart=/usr/bin/rootscript

[Install]
WantedBy=multi-user.target
